# systemd/composer-watch.service
[Unit]
Description=Composer Watch Service - Real-time document composition
Documentation=https://github.com/akaoio/composer
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/Projects/composer
ExecStart=/usr/bin/node /home/%i/Projects/composer/bin/composer-watch.js --daemon --log=/var/log/composer-watch.log --pid=/var/run/composer-watch.pid
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
Restart=on-failure
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/%i/Projects/composer/output /var/log /var/run

# Resource limits
CPUQuota=50%
MemoryLimit=512M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=composer-watch

[Install]
WantedBy=multi-user.target